# launchd/ai.bloodz.bloodzspace-deploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.bloodz.bloodzspace-deploy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bloodz/.local/bin/bloodzspace-deploy.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>StandardOutPath</key>
  <string>/Users/bloodz/Library/Logs/bloodzSpace/launchd-deploy.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/bloodz/Library/Logs/bloodzSpace/launchd-deploy.err.log</string>
</dict>
</plist>
